# supervisord.conf
[supervisord]
nodaemon=true

[program:nonebot]
command=/bin/bash /workspace/nbrun.sh

[program:lagrange]
command=/bin/bash /workspace/lgrrun.sh 